# deploy/production/cron/mosure.cron
# Install as one of:
# 1) /etc/cron.d/mosure  (recommended for server-level config)
# 2) crontab -u www-data -e

* * * * * www-data cd /var/www/mosure && /usr/bin/php artisan schedule:run >> /var/www/mosure/storage/logs/cron-schedule.log 2>&1
